# com.tars.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
╔══════════════════════════════════════════════════════════════╗
║      TARS — macOS Service (launchd)                          ║
╠══════════════════════════════════════════════════════════════╣
║  Runs TARS 24/7 as a background service.                     ║
║  Auto-restarts on crash. Starts on login.                    ║
║  Auto-reloads when Python files change.                      ║
║                                                              ║
║  Install:   ./watchdog.sh install                            ║
║  Uninstall: ./watchdog.sh uninstall                          ║
║  Status:    ./watchdog.sh status                             ║
║  Logs:      tail -f logs/tars-service.log                    ║
║             tail -f logs/tars-launchd.log                    ║
╚══════════════════════════════════════════════════════════════╝
-->
<plist version="1.0">
<dict>
    <!-- Unique identifier -->
    <key>Label</key>
    <string>com.tars.agent</string>

    <!-- Open TARS.command in Terminal.app — avoids ~/Downloads sandbox issues -->
    <!-- .command files auto-open in Terminal.app on macOS -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/open</string>
        <string>/Users/abdullah/Downloads/tars-main/TARS.command</string>
    </array>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/abdullah</string>

    <!-- Don't auto-restart the launcher — the runner handles crash recovery internally -->
    <key>KeepAlive</key>
    <false/>

    <!-- Start on user login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Throttle: wait 10s before restarting the wrapper -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Log runner stdout/stderr — /tmp avoids Downloads sandbox issues -->
    <key>StandardOutPath</key>
    <string>/tmp/tars-launcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/tars-launcher.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/abdullah/Downloads/tars-main/.venv/bin</string>
        <key>HOME</key>
        <string>/Users/abdullah</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- Process limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>
</dict>
</plist>
